# bot/web3-digest-dev.service
[Unit]
Description=Web3 Digest Bot (DEV)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/github-projects/info_denoise_compress_dev/bot
ExecStart=/home/ubuntu/github-projects/info_denoise_compress_dev/venv/bin/python main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
